# deploy/backdrop-worker.service
[Unit]
Description=Backdrop Import Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/finne/Backdrop
ExecStart=/home/finne/Backdrop/.venv/bin/python3 worker/importer.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
