# deploy/backdrop-webapp.service
[Unit]
Description=Backdrop Web App
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/finne/Backdrop
Environment=PYTHONPATH=/home/finne/Backdrop
ExecStart=/home/finne/Backdrop/.venv/bin/python3 webapp/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
